# bin/ezblock
#!/bin/bash
# /etc/init.d/ezblock

### BEGIN INIT INFO
# Provides:          sunfounder
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: ezblock service
# Description:       This service is used to startup ezblock
### END INIT INFO


case "$1" in
    start)
        echo "Starting ezblock"
        # sudo modprobe w1-gpio
        # sudo modprobe w1-therm
        sudo ezblock-service 2>&1 >> /opt/ezblock/log &
        ;;
    stop)
        echo "Stopping ezblock"
        sudo kill $(ps aux | grep -m 1 'ezblock-service' | awk '{ print $2 }')
        sudo kill $(ps aux | grep -m 1 'ezblock-service' | awk '{ print $2 }')
        sudo kill $(ps aux | grep -m 1 '/opt/ezblock/main' | awk '{ print $2 }')
        sudo kill $(ps aux | grep -m 1 '/opt/ezblock/main' | awk '{ print $2 }')
        ;;
    *)
        echo "Usage: sudo service ezblock start|stop|restart"
        exit 1
        ;;
esac

exit 0
